# bridge/systemd/kittypau-bridge.service
[Unit]
Description=Kittypau MQTT Bridge (HiveMQ -> Supabase)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kittypau
WorkingDirectory=/home/kittypau/kittypau-bridge
EnvironmentFile=/home/kittypau/kittypau-bridge/.env
ExecStart=/usr/bin/node /home/kittypau/kittypau-bridge/src/index.js

# Reinicio automatico: espera 5s entre intentos.
# Si falla mas de 5 veces en 60s, systemd lo deja caido (evita loop infinito).
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

TimeoutStopSec=10
KillSignal=SIGINT
SuccessExitStatus=0 SIGINT SIGTERM

# Logs via journald (ver con: journalctl -u kittypau-bridge -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kittypau-bridge

# Hardening minimo
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

